# .deployment/systemd/haus-facebook-refresh-access-token.timer
[Unit]
Description=Timer update haus facebook auth token

[Timer]
OnBootSec=1s
OnUnitActiveSec=2h
Persistent=true

[Install]
WantedBy=timers.target
