# .deployment/systemd/haus-jekyll-build.timer
[Unit]
Description=Run Jekyll Build every minute

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
