# .deployment/systemd/haus-facebook-info-update.timer
[Unit]
Description=Timer update haus info

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
